# clock-display.conf
[program:clock-display]
command=/usr/local/bin/clock.py
